# systemd/eink-display.service
[Unit]
Description=E-ink Display System
After=network.target
Wants=network.target

[Service]
Type=simple
User=${USER}
Group=gpio
SupplementaryGroups=spi i2c gpio dialout
WorkingDirectory=${HOME}/RpiEinky
ExecStart=${HOME}/eink_env/bin/python ${HOME}/RpiEinky/run_eink_system.py
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM
StandardOutput=append:${HOME}/RpiEinky/logs/display.log
StandardError=append:${HOME}/RpiEinky/logs/display.log

# Environment variables for hardware access
Environment=PYTHONPATH=${HOME}/RpiEinky
Environment=DISPLAY=:0
Environment=HOME=${HOME}
Environment=USER=${USER}
Environment=XDG_RUNTIME_DIR=/run/user/${UID}

# Ensure clean restarts
RemainAfterExit=no

[Install]
WantedBy=multi-user.target 