# systemd/eink-system.target
[Unit]
Description=E-ink Display System (Monitor + Upload Server)
Requires=eink-display.service eink-upload.service
After=eink-display.service eink-upload.service

[Install]
WantedBy=multi-user.target 